# django/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:gunicorn]
command=gunicorn -b 0.0.0.0:%(ENV_GUNICORN_PORT)s --reload omniport.wsgi

[program:daphne]
command=daphne -b 0.0.0.0 -p %(ENV_DAPHNE_PORT)s omniport.asgi:application